# debbuild/debian/prerm
#!/bin/sh

set -e

deb-systemd-helper disable jasper-bot.service >/dev/null
deb-systemd-helper stop jasper-bot.service >/dev/null