# debbuild/debian/prerm
#!/bin/sh

systemctl disable jasper-bot.service >/dev/null
systemctl stop jasper-bot.service >/dev/null